# basicfs/rootfs/etc/profile.d/qt4.sh
QT4DIR=/opt/qt4

pathappend /opt/qt4/bin           PATH
pathappend /opt/qt4/lib/pkgconfig PKG_CONFIG_PATH

export QT4DIR

