# basicfs/rootfs/etc/profile.d/qt4.sh
pathappend /usr/lib/qt4/bin           PATH

