# basicfs/rootfs/etc/profile.d/qt5.sh
pathappend /usr/lib/qt5/bin           PATH
pathappend /usr/lib/pkgconfig/        PKG_CONFIG_PATH
